# data/supervisor/bcron.conf
[program:bcron-sched]
directory = /etc/service
command = /usr/bin/runsv bcron-sched

[program:bcron-spool]
directory = /etc/service
command = /usr/bin/runsv bcron-spool

[program:bcron-update]
directory = /etc/service
command = /usr/bin/runsv bcron-update